# templates/init.d.mgrep-4.0.2.erb
#!/bin/bash
#
# chkconfig: 35 90 12
# description: mgrep 4.x server running on port <%= @mgrep_port %>
#
#

# Get function from functions library
. /etc/init.d/functions

MGREP_VERSION=4.0.2
PORT=<%= @mgrep_port %>
lockfile=/var/lock/subsys/mgrep-${PORT}
pidfile=/var/run/mgrep/mgrep-${PORT}.pid
pidopts="-p $pidfile"
RETVAL=0

#Defautls, can be overwritten by /etc/sysconfig/mgrep
MGREP_HOME=/usr/local/mgrep-$MGREP_VERSION
mgrepd=$MGREP_HOME/mgrep
MGREP_DICT=/srv/mgrep/mgrep-${PORT}/dict

[ -e /etc/sysconfig/mgrep ] && . /etc/sysconfig/mgrep

# Start the service mgrep
start() {
    [ -x $mgrepd ] || exit 5
    #daemon --pidfile option doesn't create pidfile for some odd reason so we are manualy creating it with the help of pgrep
    daemon --user "mgrep" "$mgrepd --port=${PORT} -f $MGREP_DICT -w $MGREP_HOME/word_divider.txt -c $MGREP_HOME/CaseFolding.txt"
    RETVAL=$?
    echo
    [ $RETVAL = 0 ] && touch ${lockfile} && pgrep -f "mgrep --port=${PORT}" > ${pidfile}
    return $RETVAL
#    success $"mgrep $MGREP_VERSION server startup"
    echo
}

# Restart the service mgrep
stop() {
    #initlog -c "echo -n Stopping mgrep $MGREP_VERSION server running on port $PORT: "
    killproc $pidopts $mgrepd
    RETVAL=$?
    echo
    [ $RETVAL = 0 ] && rm -f ${lockfile} ${pidfile}
}

rh_status() {
    status $pidopts $mgrepd
    RETVAL=$?
    return $RETVAL
}

### main logic ###
case "$1" in
  start)
        rh_status >/dev/null 2>&1 && echo "mgrep already running on port $PORT" && exit 0 
        start
        ;;
  stop)
        stop
        ;;
  status)
        #status mgrep
        rh_status
        ;;
  restart|reload|condrestart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac

exit 0
